# ansible/roles/hv-metrics-server/templates/registry-traffic-collector.timer.j2
[Unit]
Description=Collect mirror registry traffic metrics every 5 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
